# meta-alarm/recipes-extended/alarm/files/alarm.service
[Unit]
Description=Alarm service

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/alarm
